# container-files/etc/supervisor.d/selenium.conf
[program:selenium]
command=/usr/bin/java -jar /usr/bin/selenium-server-standalone.jar
stdout_logfile=/data/logs/selenium.log
stderr_logfile=/data/logs/selenium.log
priority=1010
