# container-files/etc/supervisor.d/xvfb.conf
[program:xvfb]
command=/bin/sh -c "/usr/bin/Xvfb $DISPLAY -ac -screen 0 $SCREEN_DIMENSION"
stdout_logfile=/data/logs/xvfb.log
stderr_logfile=/data/logs/xvfb.log
priority=1000
